# deploy/sniper.service
[Unit]
Description=Poly-Scout Sniper Trader
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/poly-scout
ExecStart=/usr/bin/python3 -m src.sniper_trader
Restart=always
RestartSec=5
StandardOutput=append:/var/log/sniper.log
StandardError=append:/var/log/sniper.log

# Performance tuning
Nice=-5
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50

[Install]
WantedBy=multi-user.target
