# deploy/sniper.service
[Unit]
Description=Poly-Scout Sniper Trader + Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/poly-scout
ExecStart=/bin/bash -c 'python3 -m src.sniper_trader & python3 -m src.dashboard --port 8080; wait'
Restart=always
RestartSec=5
StandardOutput=append:/var/log/sniper.log
StandardError=append:/var/log/sniper.log

# Performance tuning
Nice=-5

[Install]
WantedBy=multi-user.target
